# scripts/launchagents/com.gigaton.knowledge-extracts-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>           <string>com.gigaton.knowledge-extracts-refresh</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>cd $HOME/Documents/GitHub/master-knowledge-base &amp;&amp; /usr/bin/python3 scripts/xlsx_to_knowledge.py --output-dir knowledge-extracts/</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>          <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
  <!-- Sundays 05:00 UTC -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>       <integer>0</integer>
    <key>Hour</key>          <integer>5</integer>
    <key>Minute</key>        <integer>0</integer>
  </dict>
  <key>StandardOutPath</key> <string>/Users/admin/.knowledge_sync/knowledge-extracts.out.log</string>
  <key>StandardErrorPath</key><string>/Users/admin/.knowledge_sync/knowledge-extracts.err.log</string>
  <key>RunAtLoad</key>       <false/>
</dict>
</plist>
